# app-emulation/vmware-workstation/files/vmware.service
[Unit]
Description=VMare Workstation Linux kernel module

[Service]
ExecStart=/etc/init.d/vmware start

[Install]
WantedBy=multi-user.target
