# sys-devel/distcc/files/distccd.service
[Unit]
Description=Distributed Compiler Daemon

[Service]
Type=forking
ExecStart=/usr/sbin/distccd-wrapper

[Install]
WantedBy=multi-user.target
