# net-fs/samba/files/nmbd.service
[Unit]
Description=nmbd daemon

[Service]
ExecStart=/usr/sbin/nmbd -D
PIDFile=/var/run/nmbd.pid

[Install]
WantedBy=multi-user.target

