# net-fs/samba/files/nmbd.service
[Unit]
Description=nmbd daemon

[Service]
Type=forking
ExecStart=/usr/sbin/nmbd -D
PIDFile=/run/nmbd.pid

[Install]
WantedBy=multi-user.target

