# net-fs/samba/files/smbd.service
[Unit]
Description=smbd daemon

[Service]
Type=forking
ExecStart=/usr/sbin/smbd -D
PIDFile=/var/run/smbd.pid

[Install]
WantedBy=multi-user.target

