# net-print/cups/files/cups.service
[Unit]
Description=CUPS Server

[Service]
Type=forking
ExecStart=/usr/sbin/cupsd

[Install]
WantedBy=multi-user.target
