# dev-util/virtualbox-bin/files/virtualbox.service
[Unit]
Description=VirtualBox Linux kernel module

[Service]
ExecStart=/etc/init.d/vboxdrv start

[Install]
WantedBy=multi-user.target
